# init/systemd/gitlab.target
[Unit]
Description=GitLab - Self Hosted Git Management

[Install]
WantedBy=multi-user.target
